# salt/mlxcfg/templates/txqueuelen.rules
KERNEL=="{{ interface }}", RUN+="/sbin/ip link set {{ interface }} txqueuelen 10000"
